# deploy/systemd/fallback-backend.service
[Unit]
Description=Fallback FastAPI backend
After=network.target

[Service]
Type=simple
User=fallback
Group=fallback
WorkingDirectory=/opt/fallback/backend
ExecStart=/opt/fallback/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
